# deployments/systemd/plantx-order.service
[Unit]
Description=PlantX Order Service
After=network.target

[Service]
Type=simple
User=plantx
Group=plantx
ExecStart=/opt/plantx/bin/order-service
Restart=always
RestartSec=5
Environment="GRPC_PORT=8080"
Environment="HTTP_PORT=8081"
Environment="DATABASE_DSN=postgres://plantx:plantx@localhost:5432/plantx?sslmode=disable"
Environment="NATS_URL=nats://localhost:4222"
Environment="MAXKEY_ISSUER=http://localhost:8080"
Environment="OPA_URL=http://localhost:8181"

[Install]
WantedBy=multi-user.target
